# services/depacketizer.service
[Unit]
Description=Depacketizer

[Service]
Type=simple
ExecStart=/bin/python -u /root/bus/depacketizer.py