# services/echo.service
[Unit]
Description=Packet Echo

[Service]
Type=simple
ExecStart=/bin/python -u /root/bus/echo.py
