# services/flash-fpga.service
[Unit]
Description=Flash FPGA
Before=fpga.service

[Service]
Type=oneshot
ExecStart=/bin/sleep 1
ExecStart=/bin/python -u /root/fpga/flash.py

[Install]
WantedBy=default.target
