# services/housekeeping.service
[Unit]
Description=Housekeeping

[Service]
Type=simple
ExecStart=/bin/python -u /root/housekeeping/housekeeping.py