# services/pat.service
[Unit]
Description=PAT

[Service]
Type=simple
ExecStart=/bin/pat