# services/pat.service
[Unit]
Description=PAT

[Service]
Type=simple
StandardOutput=null
ExecStart=/root/bin/pat
